# systemd/trustpanel-cleanup.service
[Unit]
Description=TrustPanel Cleanup

[Service]
Type=oneshot
ExecStart=/opt/trustpanel/venv/bin/python -u /opt/trustpanel/services/cleanup.py
StandardOutput=journal
StandardError=journal
