# systemd/trustpanel-backup.service
[Unit]
Description=TrustPanel Backup

[Service]
Type=oneshot
ExecStart=/opt/trustpanel/venv/bin/python /opt/trustpanel/services/backup.py
