# systemd/trustpanel-bot.service
[Unit]
Description=TrustPanel Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/opt/trustpanel
EnvironmentFile=/opt/trustpanel/.env
ExecStart=/opt/trustpanel/venv/bin/python /opt/trustpanel/bot/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
